# .config/systemd/user/mailsync.service
[Unit]
Description=Sync mailboxes

[Service]
Type=simple
ExecStart=/home/tuomo/.local/bin/mailsync

[Install]
WantedBy=default.target
